# example_service/vuln/templates/example_service.service
[Unit]
Description=Example service starter

[Service]
WorkingDirectory={{ challenge_path }}/
ExecStart=/usr/bin/python3 example_service.py
